# etc/conf/supervisor.conf
[program:ej-server]
command=inv run-production
autostart=true
autorestart=true
env=DJANGO_DEBUG=false
stderr_logfile=/var/log/ej-server.err.log
stdout_logfile=/var/log/ej-server.out.log
